# deploy/launchd/io.hotchkiss.recon-demo.refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  AE.5 — launchd plist for the nightly refresh job. Fires at 03:00
  America/New_York every day. Rebuilds both demo SQLite databases
  (spec + sasquatch) from scratch using the latest recon-gen PyPI
  wheel, then kickstart-restarts both server jobs (~5s outage each).

  Install:
    cp deploy/launchd/io.hotchkiss.recon-demo.refresh.plist \
       ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/io.hotchkiss.recon-demo.refresh.plist

  The plist fires the wrapper script at /Users/recon-demo/bin/refresh-demos.sh
  (operator installs from deploy/launchd/refresh-demos.sh during AE.1).
  StartCalendarInterval in launchd is local-timezone — set the Mac
  mini's system timezone to America/New_York or adjust the Hour field.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.hotchkiss.recon-demo.refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/recon-demo/bin/refresh-demos.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/recon-demo</string>

    <!--
      03:00 every day, local time. launchd interprets StartCalendarInterval
      against the system's local timezone (`systemsetup -gettimezone`).
      If the Mac mini sits on UTC, change Hour to 7 (03:00 ET in winter)
      or 8 (03:00 ET in summer) — daylight-savings drift is fine for
      a demo's refresh schedule.
    -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!--
      RunAtLoad false: a fresh `launchctl load` shouldn't trigger an
      immediate refresh (the server jobs handle that on their own
      RunAtLoad=true). Refresh fires only at 03:00.
    -->
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/recon-demo/logs/refresh.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/recon-demo/logs/refresh.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/recon-demo/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/recon-demo</string>
    </dict>
</dict>
</plist>
