# deploy/launchd/io.hotchkiss.recon-demo.sasquatch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  AE.4 — launchd plist for the sasquatch.recon-gen.hotchkiss.io demo
  (recon-gen studio --demo-mode against sasquatch_pr L2, sandboxed
  with mutation perimeter cut).

  Install:
    cp deploy/launchd/io.hotchkiss.recon-demo.sasquatch.plist \
       ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/io.hotchkiss.recon-demo.sasquatch.plist

  Reload after editing:
    launchctl unload -w ~/Library/LaunchAgents/io.hotchkiss.recon-demo.sasquatch.plist
    launchctl load   -w ~/Library/LaunchAgents/io.hotchkiss.recon-demo.sasquatch.plist

  The plist invokes `recon-gen studio --demo-mode`. AE.2.b's
  --demo-mode flag strips POST/PUT/DELETE /l2_shape/* + PUT
  /data/knobs/etl_hook + POST /deploy from the route table. Trainer
  knob mutations (plants, end_date, window, seed, scope) land in a
  per-process tmpdir for .studio-state.yaml — wiped on every server
  restart (nightly refresh restart cycles state freshly).

  The sandbox-exec profile (recon-demo-sasquatch.sb) is defense in
  depth: it denies file-write on the L2 yaml + cfg.yaml regardless
  of UI state. STUDIO_STATE_DIR points at the per-process tmpdir;
  the wrapper script (set up by AE.1 / the install step) computes a
  fresh tmpdir per launchd start.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.hotchkiss.recon-demo.sasquatch</string>

    <!--
      The wrapper-shell pattern. launchd doesn't support `$(...)` in
      ProgramArguments, so a thin /bin/sh -c wrapper is the cleanest
      way to compute a fresh STUDIO_STATE_DIR per launchd start
      (mktemp -d) before exec-ing into sandbox-exec.

      Wrapper pseudocode:
        STUDIO_STATE_DIR="$(mktemp -d -t recon-demo-studio-state)"
        exec sandbox-exec -D HOME=... -D INSTANCE_DIR=... \
             -D PORT=8402 -D PYTHON=... \
             -D STUDIO_STATE_DIR="$STUDIO_STATE_DIR" \
             -f /Users/recon-demo/sandbox/recon-demo-sasquatch.sb \
             -- /Users/recon-demo/venv/bin/recon-gen studio --demo-mode ...

      Operator authors this wrapper at /Users/recon-demo/bin/launch-sasquatch.sh
      during AE.1. It's owned by recon-demo, mode 0500 (read+exec, no write),
      so the sandboxed process can't tamper with it.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/recon-demo/bin/launch-sasquatch.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/recon-demo/sasquatch_pr</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/recon-demo/sasquatch_pr/logs/server.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/recon-demo/sasquatch_pr/logs/server.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/recon-demo/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/recon-demo</string>
    </dict>
</dict>
</plist>
