# deploy/launchd/io.hotchkiss.recon-demo.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  AE.6 — LaunchDaemon plist for the Cloudflare Tunnel daemon
  (cloudflared). Runs the tunnel as a long-lived process under
  recon-demo, bridging the two locally-bound dashboard servers to
  their public hostnames via the operator's Cloudflare account.

  LaunchDaemon (not Agent) for the same Tahoe reason as the other
  recon-demo daemons. UserName drops privileges to recon-demo.

  Prerequisite (AE.6 install steps):
    1. brew install cloudflared (admin)
    2. As recon-demo: cloudflared tunnel login (browser auth)
    3. As recon-demo: cloudflared tunnel create recon-demo
    4. As recon-demo: copy deploy/cloudflared/cloudflared.yml.example
       to /Users/recon-demo/.cloudflared/config.yml, substitute UUID
    5. As recon-demo: cloudflared tunnel route dns recon-demo <host>
    6. Install this plist + bootstrap.

  Path differs on Apple Silicon (/opt/homebrew/bin) vs Intel
  (/usr/local/bin). Verify with `which cloudflared` and edit if needed.

  Install (as admin):
    sudo cp deploy/launchd/io.hotchkiss.recon-demo.tunnel.plist \
            /Library/LaunchDaemons/
    sudo chown root:wheel /Library/LaunchDaemons/io.hotchkiss.recon-demo.tunnel.plist
    sudo chmod 0644 /Library/LaunchDaemons/io.hotchkiss.recon-demo.tunnel.plist
    sudo launchctl bootstrap system \
        /Library/LaunchDaemons/io.hotchkiss.recon-demo.tunnel.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.hotchkiss.recon-demo.tunnel</string>

    <key>UserName</key>
    <string>recon-demo</string>
    <key>GroupName</key>
    <string>staff</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/cloudflared</string>
        <string>tunnel</string>
        <string>--config</string>
        <string>/Users/recon-demo/.cloudflared/config.yml</string>
        <string>run</string>
        <string>recon-demo</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/recon-demo</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/recon-demo/logs/tunnel.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/recon-demo/logs/tunnel.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/recon-demo</string>
    </dict>
</dict>
</plist>
