# deploy/launchd/io.hotchkiss.recon-demo.refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  AE.5 — LaunchDaemon plist for the nightly refresh job. Fires at
  03:00 local time daily. Rebuilds both demo SQLite databases from
  scratch using the latest recon-gen PyPI wheel, then restarts both
  server daemons (~5s outage each).

  LaunchDaemon (not Agent) for the same Tahoe reason as the server
  daemons. UserName drops privileges to recon-demo before running
  refresh-demos.sh.

  Install (as admin):
    sudo cp deploy/launchd/io.hotchkiss.recon-demo.refresh.plist \
            /Library/LaunchDaemons/
    sudo chown root:wheel /Library/LaunchDaemons/io.hotchkiss.recon-demo.refresh.plist
    sudo chmod 0644 /Library/LaunchDaemons/io.hotchkiss.recon-demo.refresh.plist
    sudo launchctl bootstrap system \
        /Library/LaunchDaemons/io.hotchkiss.recon-demo.refresh.plist

  StartCalendarInterval Hour/Minute is local-time. Operator running
  on a UTC-system Mac mini should adjust Hour to match their target
  local-time (e.g. 7 for 03:00 EST during winter, 8 for EDT).
  Daylight savings drift across the year is fine for a demo
  refresh — it'll slide by one hour twice yearly.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.hotchkiss.recon-demo.refresh</string>

    <key>UserName</key>
    <string>recon-demo</string>
    <key>GroupName</key>
    <string>staff</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/recon-demo/bin/refresh-demos.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/recon-demo</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/recon-demo/logs/refresh.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/recon-demo/logs/refresh.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/recon-demo/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/recon-demo</string>
    </dict>
</dict>
</plist>
