# deploy/launchd/io.hotchkiss.recon-demo.spec.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  AE.4 — LaunchDaemon plist for the spec.recon-gen.hotchkiss.io demo
  (recon-gen dashboards against spec_example L2, sandboxed).

  This is a LaunchDaemon, NOT a LaunchAgent. On macOS 26 (Tahoe) a
  hidden user without an active login session can't host LaunchAgents
  in the per-user domain — `launchctl bootstrap user/<uid> ...` returns
  "Input/output error" even when the plist is valid XML and all
  ProgramArguments paths exist. LaunchDaemons run in the system domain
  and drop privileges via UserName/GroupName at process start, which
  works for our recon-demo-without-GUI-session shape.

  Install (as admin):
    sudo cp deploy/launchd/io.hotchkiss.recon-demo.spec.plist \
            /Library/LaunchDaemons/
    sudo chown root:wheel /Library/LaunchDaemons/io.hotchkiss.recon-demo.spec.plist
    sudo chmod 0644 /Library/LaunchDaemons/io.hotchkiss.recon-demo.spec.plist
    sudo launchctl bootstrap system \
        /Library/LaunchDaemons/io.hotchkiss.recon-demo.spec.plist

  Restart after editing:
    sudo launchctl bootout system \
        /Library/LaunchDaemons/io.hotchkiss.recon-demo.spec.plist
    sudo launchctl bootstrap system \
        /Library/LaunchDaemons/io.hotchkiss.recon-demo.spec.plist

  Status: `sudo launchctl print system/io.hotchkiss.recon-demo.spec`
  Logs:   /Users/recon-demo/spec_example/logs/server.{out,err}.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.hotchkiss.recon-demo.spec</string>

    <!-- Drop privileges to recon-demo at process start. -->
    <key>UserName</key>
    <string>recon-demo</string>
    <key>GroupName</key>
    <string>staff</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/sandbox-exec</string>
        <string>-D</string>
        <string>HOME=/Users/recon-demo</string>
        <string>-D</string>
        <string>INSTANCE_DIR=/Users/recon-demo/spec_example</string>
        <string>-D</string>
        <string>PORT=8401</string>
        <string>-D</string>
        <string>PYTHON=/Users/recon-demo/venv/bin/python3.13</string>
        <string>-f</string>
        <string>/Users/recon-demo/sandbox/recon-demo-spec.sb</string>
        <string>--</string>
        <string>/Users/recon-demo/venv/bin/recon-gen</string>
        <string>dashboards</string>
        <string>-c</string>
        <string>/Users/recon-demo/spec_example/config.yaml</string>
        <string>--l2</string>
        <string>/Users/recon-demo/spec_example/l2.yaml</string>
        <string>--port</string>
        <string>8401</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--no-docs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/recon-demo/spec_example</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/recon-demo/spec_example/logs/server.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/recon-demo/spec_example/logs/server.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/recon-demo/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/recon-demo</string>
        <!--
          Disable Python bytecode caching. Default behavior tries to
          write .pyc files to ``<sys.path-entry>/__pycache__/``, which
          includes the system Python stdlib at
          ``/opt/homebrew/Cellar/python@3.13/.../lib/python3.13/...``.
          The sandbox profile denies file-write outside the per-instance
          state dir — Python normally swallows the EPERM but uvicorn
          surfaces it during its accept-loop setup as
          ``ERROR: [Errno 1] Operation not permitted`` and crashes the
          server. PYTHONDONTWRITEBYTECODE=1 makes Python skip the .pyc
          write attempts entirely; modest import-startup penalty for a
          long-running server.
        -->
        <key>PYTHONDONTWRITEBYTECODE</key>
        <string>1</string>
    </dict>
</dict>
</plist>
