# deploy/launchd/io.hotchkiss.recon-demo.sasquatch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  AE.4 — LaunchDaemon plist for the sasquatch.recon-gen.hotchkiss.io
  demo (recon-gen studio --demo-mode against sasquatch_pr L2,
  sandboxed with mutation perimeter cut).

  LaunchDaemon (not Agent) for the same Tahoe reason as the spec
  daemon — hidden user, no GUI session, system-domain load is what
  works. UserName/GroupName drops privileges to recon-demo at process
  start.

  Install (as admin):
    sudo cp deploy/launchd/io.hotchkiss.recon-demo.sasquatch.plist \
            /Library/LaunchDaemons/
    sudo chown root:wheel /Library/LaunchDaemons/io.hotchkiss.recon-demo.sasquatch.plist
    sudo chmod 0644 /Library/LaunchDaemons/io.hotchkiss.recon-demo.sasquatch.plist
    sudo launchctl bootstrap system \
        /Library/LaunchDaemons/io.hotchkiss.recon-demo.sasquatch.plist

  The plist execs launch-sasquatch.sh (in /Users/recon-demo/bin/).
  That wrapper computes STUDIO_STATE_DIR via mktemp -d (since
  launchd's ProgramArguments doesn't expand $(...)) before
  exec-ing sandbox-exec + `recon-gen studio --demo-mode`.

  Requires recon-gen >=11.6.5 on PyPI (--demo-mode shipped there).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.hotchkiss.recon-demo.sasquatch</string>

    <key>UserName</key>
    <string>recon-demo</string>
    <key>GroupName</key>
    <string>staff</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/recon-demo/bin/launch-sasquatch.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/recon-demo/sasquatch_pr</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/recon-demo/sasquatch_pr/logs/server.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/recon-demo/sasquatch_pr/logs/server.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/recon-demo/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/recon-demo</string>
        <!--
          See spec plist for the rationale — sandbox-exec denies .pyc
          writes to system Python stdlib paths; uvicorn surfaces the
          downstream EPERM as a crash. PYTHONDONTWRITEBYTECODE=1
          skips bytecode caching entirely.
        -->
        <key>PYTHONDONTWRITEBYTECODE</key>
        <string>1</string>
    </dict>
</dict>
</plist>
